# d800_download.service
[Unit]
Description=Téléchargement automatique des photos depuis Nikon D800
After=network.target

[Service]
Type=simple
User=server01
WorkingDirectory=/home/server01/projet_ftp/Projet_FTP
ExecStart=/usr/bin/python3 /home/server01/projet_ftp/Projet_FTP/d800_auto_download.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
