# photo-ftp-web.service
[Unit]
Description=Service web de transfert de photos
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/projet_ftp/Projet_FTP
ExecStart=/usr/bin/python3 simple_webui.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=photo-ftp-web

[Install]
WantedBy=multi-user.target
