# photo-ftp-web.service
[Unit]
Description=Photo Transfer Web UI - Configuration interface for photo transfer service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/photo-ftp
ExecStart=/usr/bin/python3 /home/pi/photo-ftp/webui.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Variables d'environnement
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi/photo-ftp

[Install]
WantedBy=multi-user.target
